# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading
autorestart=true
priority=10
stdout_logfile=/var/log/pulseaudio.log
stderr_logfile=/var/log/pulseaudio.err

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 1920x1080x24 -ac +extension GLX +extension RANDR +render -noreset
autorestart=true
priority=20
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:xrandr-setup]
command=/opt/scripts/setup-display.sh
environment=DISPLAY="%(ENV_DISPLAY)s"
autorestart=false
startsecs=0
exitcodes=0
priority=25
stdout_logfile=/var/log/xrandr-setup.log
stderr_logfile=/var/log/xrandr-setup.err

[program:xfce]
command=startxfce4
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer"
autorestart=true
priority=30
startsecs=5
startretries=3
stdout_logfile=/var/log/xfce.log
stderr_logfile=/var/log/xfce.err

[program:selkies-gstreamer]
command=selkies-gstreamer --addr=0.0.0.0 --port=6901 --enable_https=false --encoder=nvh264enc --framerate=60 --video_bitrate=6000 --enable_resize=false --enable_audio=true --web_root=/opt/gst-web
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer",GST_DEBUG="2",PULSE_SERVER="/var/run/pulse/native",GSTREAMER_PATH="/opt/gstreamer",PATH="/opt/gstreamer/bin:/opt/VirtualGL/bin:/usr/games:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",LD_LIBRARY_PATH="/opt/gstreamer/lib/x86_64-linux-gnu",GST_PLUGIN_PATH="/opt/gstreamer/lib/x86_64-linux-gnu/gstreamer-1.0",GST_PLUGIN_SYSTEM_PATH="/opt/gstreamer/lib/x86_64-linux-gnu/gstreamer-1.0:/usr/lib/x86_64-linux-gnu/gstreamer-1.0",GI_TYPELIB_PATH="/opt/gstreamer/lib/x86_64-linux-gnu/girepository-1.0:/usr/lib/x86_64-linux-gnu/girepository-1.0",PYTHONPATH="/opt/gstreamer/lib/python3/dist-packages"
autorestart=true
priority=50
startsecs=5
startretries=5
stdout_logfile=/var/log/selkies-gstreamer.log
stderr_logfile=/var/log/selkies-gstreamer.err

[program:steam]
command=/opt/VirtualGL/bin/vglrun /usr/games/steam -no-browser +open steam://open/minigameslist
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer",PATH="/opt/VirtualGL/bin:/usr/games:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",VGL_DISPLAY=egl
autorestart=true
priority=60
startsecs=10
startretries=3
stdout_logfile=/var/log/steam.log
stderr_logfile=/var/log/steam.err
