# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading
autorestart=true
priority=10
stdout_logfile=/var/log/pulseaudio.log
stderr_logfile=/var/log/pulseaudio.err

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 1920x1080x24 -ac +extension GLX +extension RANDR +render -noreset
autorestart=true
priority=20
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:xrandr-setup]
command=/opt/scripts/setup-display.sh
environment=DISPLAY="%(ENV_DISPLAY)s"
autorestart=false
startsecs=0
exitcodes=0
priority=25
stdout_logfile=/var/log/xrandr-setup.log
stderr_logfile=/var/log/xrandr-setup.err

[program:xfce]
command=startxfce4
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer"
autorestart=true
priority=30
startsecs=5
startretries=3
stdout_logfile=/var/log/xfce.log
stderr_logfile=/var/log/xfce.err

[program:selkies-gstreamer]
command=/opt/scripts/start-selkies.sh
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer",GST_DEBUG="2",PULSE_SERVER="/var/run/pulse/native"
autorestart=true
priority=50
startsecs=5
startretries=5
stdout_logfile=/var/log/selkies-gstreamer.log
stderr_logfile=/var/log/selkies-gstreamer.err

[program:steam]
command=/opt/VirtualGL/bin/vglrun /usr/games/steam -no-browser -no-cef-sandbox +open steam://open/minigameslist
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer",PATH="/opt/VirtualGL/bin:/usr/games:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin",VGL_DISPLAY=egl,STEAM_DISABLE_BROWSER_QT_DATAPATH=1,STEAMWEBHELPER_EXTERNAL=1
autorestart=true
priority=60
startsecs=10
startretries=3
stdout_logfile=/var/log/steam.log
stderr_logfile=/var/log/steam.err
