# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log

[program:xvfb]
command=Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_RESOLUTION)s -ac +extension GLX +render -noreset
environment=DISPLAY="%(ENV_DISPLAY)s"
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:xfce]
command=startxfce4
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer",LIBGL_ALWAYS_SOFTWARE="1",GALLIUM_DRIVER="llvmpipe"
autorestart=true
priority=20
startsecs=5
startretries=3
stdout_logfile=/var/log/xfce.log
stderr_logfile=/var/log/xfce.err

[program:kasmvnc]
command=kasmvncserver %(ENV_DISPLAY)s -websocketPort 6901 -fg -disableBasicAuth
user=gamer
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="/home/gamer"
autorestart=true
priority=30
startsecs=3
stdout_logfile=/var/log/kasmvnc.log
stderr_logfile=/var/log/kasmvnc.err

[program:pulseaudio]
command=pulseaudio --system --disallow-exit --disallow-module-loading
autorestart=true
priority=15
stdout_logfile=/var/log/pulseaudio.log
stderr_logfile=/var/log/pulseaudio.err
